# faveo/supervisor.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
user = root
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[supervisorctl]

[program:ssh]
command=/usr/sbin/sshd -e -D
autostart=true
autorestart=true
priority=995
stopsignal=INT
stopasgroup=true
killasgroup=true

[program:mysql]
command=su - mysql -s /bin/sh -c "/usr/bin/mysqld_safe &"
autorestart=true
user=root
priority=996

[program:php-fpm]
command=/usr/sbin/php-fpm7.1 -F
autostart=true
autorestart=true
priority=997
stdout_logfile=/var/log/supervisor/php-fpm/stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/php-fpm/stderr.log
stderr_logfile_maxbytes=0
exitcodes=0


[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
priority=998
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/nginx/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/nginx/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB

[program:cron]
command=cron -f
autostart=true
autorestart=true
priority=999
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/cron/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/cron/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB